# deploying_tools/ccv.conf
[program:ccv]
command=/home/nanoafrika/ccv_supervisior
user=nanoafrika
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/nanoafrika/logs/ccv.log